# infrastructure/systemd/genesis-mesh-node-d.service
[Unit]
Description=Genesis Mesh Node D (backup router)
After=genesis-mesh-na.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
User=azureuser
WorkingDirectory=/opt/genesis-mesh
Environment=GENESIS_LOG_LEVEL=INFO
Environment=GENESIS_LOG_FORMAT=text
ExecStart=/opt/genesis-mesh/.venv/bin/genesis-mesh join --na https://na.genesismesh.connectorzzz.com --config /home/azureuser/.genesis-mesh-node-d/config.toml --persistent --listen-port 7444
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
